# Loading the system/sources files/etc/inittab
::sysinit:/etc/init.d/rc.S
console::respawn:-/bin/ash
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::halt:/bin/umount -a -r
